# clouddeploy/playbooks/roles/host_setup/templates/cpu_isolation.service.j2
[Unit]
Description=Setting the CPU affinity for kthreads and IRQs

[Service]
Type=oneshot
ExecStart=/bin/sleep 10
ExecStart=/usr/local/bin/cpu_isolation.sh

[Install]
WantedBy=multi-user.target
