# clouddeploy/playbooks/roles/haproxy_server/files/supervisord.conf
[supervisord]
nodaemon=true

[program:rsyslogd]
command=/usr/sbin/rsyslogd -n -c5
priority=1
autorestart=true
redirect_stderr=true

[program:keepalived]
command=bash -c "/docker/haproxy/ka_clean.sh;exec /usr/sbin/keepalived -ndlD -p /run/keepalived.pid"
priority=2
autorestart=true

[program:haproxy]
command=haproxy -f /etc/haproxy/haproxy.cfg -db
stderr_logfile=/var/log/haproxy/haproxy.errors
autorestart=true

